# .config/autostart/gimp is sick.desktop
[Desktop Entry]
Exec=gimp
Name=gimp is sick
Type=Application
Version=1.0
X-LXQt-Autostart-disabled=true
X-LXQt-X11-Only=true
