# .config/autostart/sxhkd startup.desktop
[Desktop Entry]
Exec=sxhkd
Name=sxhkd startup
Type=Application
Version=1.0
X-LXQt-X11-Only=true
