# .config/autostart/sxhkd startup.desktop
[Desktop Entry]
Exec=killall -q sxhkd; sxhkd &
Name=sxhkd startup
Type=Application
Version=1.0
