# .config/autostart/sxhkd startup.desktop
[Desktop Entry]
Exec=killall -q sxhkd; sxhkd &
Name=sxhkd startup
Type=Application
Version=1.0
X-LXQt-Autostart-disabled=true
X-LXQt-X11-Only=true
