# .config/autostart/warpd for vimium on computer.desktop
[Desktop Entry]
Exec=warpd
Name=warpd for vimium on computer
Type=Application
Version=1.0
