# .config/autostart/XScreenSaver.desktop
[Desktop Entry]
Exec=xscreensaver -no-splash
Name=XScreenSaver
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Autostart-disabled=true
